# hoomenets/files/etc/init.d/hmnet
#!/bin/sh /etc/rc.common

START=100

boot() {
	start
}

restart() {
	stop
	start
}

start() {
	hmnetup &
}

stop() {
	hmnetdown
}